# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command=/usr/sbin/apache2ctl -DFOREGROUND
autorestart=true
startretries=5